# scripts/com.sasha.dotfiles.autocommit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sasha.dotfiles.autocommit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>cd "$HOME/dotfiles" &amp;&amp; git add -A &amp;&amp; { git diff --cached --quiet 2>/dev/null || git commit -m auto 2>/dev/null &amp;&amp; git push 2>/dev/null; }</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
</dict>
</plist>
